# src/device/BMP180PollerDaemon.conf
[program:BMP180PollerDaemon]
command=/usr/bin/python /home/odroid/Adafruit_Python_BMP/examples/bmp180_mysql.py
#directory==/home/odroid/Adafruit_Python_BMP/examples/
stdout_logfile=/home/odroid/Adafruit_Python_BMP/examples/app_sd.log
autorestart = true
redirect_stderr=true
stdout_events_enabled=false
loglevel=warn
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
